# deploy/systemd/bookreserch-server.service
[Unit]
Description=BookReserch Fastify Server
After=network.target

[Service]
Type=simple
User=sachun78
WorkingDirectory=/opt/bookReserch
Environment=NODE_ENV=production
EnvironmentFile=-/opt/bookReserch/apps/server/.env
ExecStart=/usr/bin/npm run start -w @bookreserch/server
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
