# deploy/systemd/library-insights-server.service
[Unit]
Description=LibraryInsights Fastify Server
After=network.target

[Service]
Type=simple
User=sachun78
WorkingDirectory=/opt/library-insights
Environment=NODE_ENV=production
EnvironmentFile=-/opt/library-insights/apps/server/.env
ExecStart=/usr/bin/npm run start -w @library-insights/server
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
